# scripts/launchd/com.cockpit.structural-levels.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cockpit.structural-levels</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ravikumarbasker/trading-cockpit/scripts/launchd/structural-levels.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ravikumarbasker/trading-cockpit</string>

  <!-- Weekdays (Mon-Fri = 1-5) at 9:23 AM local time, pre-RTH. -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>23</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/ravikumarbasker/Library/Logs/cockpit-structural-levels.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ravikumarbasker/Library/Logs/cockpit-structural-levels.stderr.log</string>
</dict>
</plist>
