# scripts/launchd/com.cockpit.mbo-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cockpit.mbo-ingest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ravikumarbasker/trading-cockpit/scripts/launchd/mbo-ingest.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/ravikumarbasker/trading-cockpit</string>

  <!-- Every hour at :23 (matches the Claude session-only cron e7cd4e33) -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>23</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/ravikumarbasker/Library/Logs/cockpit-mbo-ingest.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ravikumarbasker/Library/Logs/cockpit-mbo-ingest.stderr.log</string>
</dict>
</plist>
